# launchagents/obsidian-nightly-knowledge-review.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>__LABEL_NS__.obsidian-nightly-knowledge-review</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.local/bin/obsidian-nightly-knowledge-review.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>__NIGHTLY_HOUR__</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/obsidian-nightly-knowledge-review.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/obsidian-nightly-knowledge-review.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>OBSIDIAN_VAULT</key>
        <string>__VAULT__</string>
    </dict>
</dict>
</plist>
